# recipes-images/base-files/base-files/adduser.service
[Unit]
Description=add user ... workaround for broken useradd.bbclass
Before=xserver-nodm.service

[Service]
ExecStart=/usr/bin/adduser.sh


